# squashfs-root/etc/init.d/button
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=10

start() {
    mknod /dev/gpio c 252 0
    button i 1 &
	button l &
}
