# squashfs-root/etc/init.d/guest_control
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=97

start() {
	guest_control &
}
stop() {
	killall guest_control &
}
