# squashfs-root/usr/lib/opkg/info/dev_mng.list
/etc/oui.txt
/usr/bin/dev_mng
/usr/bin/ip_eventd
/etc/init.d/dev_mng
