# squashfs-root/usr/lib/opkg/info/guest_control.list
/usr/bin/guest_control
/etc/init.d/guest_control
