# squashfs-root/usr/lib/opkg/info/luci-app-parentctl.list
/usr/lib/lua/luci/i18n/parentctl.zh-cn.lmo
/etc/config/parentctl
/usr/lib/lua/luci/view/parentctl/parentctl.htm
/etc/init.d/parentctl
/usr/lib/lua/luci/controller/admin/parentctl.lua
