# squashfs-root/usr/lib/opkg/info/luci-app-safeset.list
/usr/lib/lua/luci/view/safeset/safeset.htm
/etc/config/safeset
/etc/init.d/safeset
/usr/lib/lua/luci/controller/safeset.lua
