# squashfs-root/usr/lib/opkg/info/tr069.list
/etc/init.d/tr069
/usr/bin/http_event_inform
