# squashfs-root/usr/bin/device_manage_set
#!/bin/sh /etc/rc.common
# Parameter: MAC IP if_type block_user tx_rate rx_rate;

CONFIG="device_manage"
TYPE="limit"
ARP_SET="/sbin/arp"
MAC="$1"
IP="$2"
if_type="$3"
block_user="$4"
tx_rate="$5"
rx_rate="$6"
section_name=${MAC//:/_}

uci set wireless.mt7620.ac_policy=2
uci set wireless.mt7612e.ac_policy=2
iwpriv ra0 set AccessPolicy=2
iwpriv rai0 set AccessPolicy=2

set_rx_rate() {
	local action="$1"
	local drate="$2"
	if [ "$action" = "add" ]; then
		iptables -t mangle -w -A limit_chain --dst $IP -m hashlimit --hashlimit-name dst_$section_name \
			--hashlimit ${drate}kb/s --hashlimit-burst ${drate}kb/s --hashlimit-mode dstip -j RETURN >/dev/null 2>&1
		iptables -t mangle -w -A limit_chain --dst $IP -j DROP >/dev/null 2>&1
	elif [ "$action" = "del" ]; then
		iptables -t mangle -w -D limit_chain --dst $IP -m hashlimit --hashlimit-name dst_$section_name \
			--hashlimit ${drate}kb/s --hashlimit-burst ${drate}kb/s --hashlimit-mode dstip -j RETURN >/dev/null 2>&1
		iptables -t mangle -w -D limit_chain --dst $IP -j DROP >/dev/null 2>&1
	fi
}

set_tx_rate() {
	local action="$1"
	local urate="$2"
	if [ "$action" = "add" ]; then
		iptables -t mangle -w -A limit_chain --src $IP -m hashlimit --hashlimit-name src_$section_name \
			--hashlimit ${urate}kb/s --hashlimit-burst ${urate}kb/s --hashlimit-mode srcip -j RETURN >/dev/null 2>&1
		iptables -t mangle -w -A limit_chain --src $IP -j DROP >/dev/null 2>&1
	elif [ "$action" = "del" ]; then
		iptables -t mangle -w -D limit_chain --src $IP -m hashlimit --hashlimit-name src_$section_name \
			--hashlimit ${urate}kb/s --hashlimit-burst ${urate}kb/s --hashlimit-mode srcip -j RETURN >/dev/null 2>&1
		iptables -t mangle -w -D limit_chain --src $IP -j DROP >/dev/null 2>&1
	fi
}

if [ "$TYPE" = "$(uci get $CONFIG.$section_name)" ]; then
	block_user_old="$(uci get $CONFIG.$section_name.block_user)"
	tx_rate_old="$(uci get $CONFIG.$section_name.tx_rate)"
	rx_rate_old="$(uci get $CONFIG.$section_name.rx_rate)"

	if [ $block_user_old -eq 0 -a $block_user -ne 0 ]; then
		uci set $CONFIG.$section_name.block_user="$block_user"
		if [ "$if_type" = "0" ]; then
			arptables -A INPUT --src-mac $MAC -j DROP
			$ARP_SET -d $IP
		else
			iwpriv ra0 set ACLAddEntry="$MAC"
			iwpriv rai0 set ACLAddEntry="$MAC"
			uci add_list wireless.mt7620.ctrl_list="$MAC"
			uci add_list wireless.mt7612e.ctrl_list="$MAC"
		fi
	fi

	if [ $block_user_old -ne 0 -a $block_user -eq 0 ]; then
		uci set $CONFIG.$section_name.block_user="$block_user"
		if [ "$if_type" = "0" ]; then
			arptables -D INPUT --src-mac $MAC -j DROP
		else
			iwpriv ra0 set ACLDelEntry="$MAC"
			iwpriv rai0 set ACLDelEntry="$MAC"
			uci del_list wireless.mt7620.ctrl_list="$MAC"
			uci del_list wireless.mt7612e.ctrl_list="$MAC"
		fi
	fi

	if [ $tx_rate -ne $tx_rate_old -a $block_user_old -eq 0 ]; then
		uci set $CONFIG.$section_name.tx_rate="$tx_rate"
		if [ $tx_rate_old -gt 0 ]; then
			set_tx_rate "del" $tx_rate_old
		fi
		if [ $tx_rate -gt 0 ]; then
			set_tx_rate "add" $tx_rate
		fi
	fi

	if [ $rx_rate -ne $rx_rate_old -a $block_user_old -eq 0 ]; then
		uci set $CONFIG.$section_name.rx_rate="$rx_rate"
		if [ $rx_rate_old -gt 0 ]; then
			set_rx_rate "del" $rx_rate_old
		fi
		if [ $rx_rate -gt 0 ]; then
			set_rx_rate "add" $rx_rate
		fi
	fi

	if [ $block_user -eq 0 -a  $tx_rate -eq 0 -a $rx_rate -eq 0 ]; then
		uci delete $CONFIG.$section_name
	fi
elif [ $block_user -ne 0 -o  $tx_rate -gt 0 -o $rx_rate -gt 0 ]; then
	uci set $CONFIG.$section_name="$TYPE"
	uci set $CONFIG.$section_name.mac="$MAC"
	uci set $CONFIG.$section_name.if_type="$if_type"
	uci set $CONFIG.$section_name.block_user="$block_user"
	uci set $CONFIG.$section_name.tx_rate="$tx_rate"
	uci set $CONFIG.$section_name.rx_rate="$rx_rate"
	if [ $block_user -ne 0 ]; then
		if [ "$if_type" = "0" ]; then
			arptables -A INPUT --src-mac $MAC -j DROP
			$ARP_SET -d $IP
		else
			uci add_list wireless.mt7620.ctrl_list="$MAC"
			uci add_list wireless.mt7612e.ctrl_list="$MAC"
			iwpriv ra0 set ACLAddEntry="$MAC"
			iwpriv rai0 set ACLAddEntry="$MAC"
		fi
	fi
	if [ $tx_rate -gt 0 ]; then
		set_tx_rate "add" $tx_rate
	fi
	if [ $rx_rate -gt 0 ]; then
		set_rx_rate "add" $rx_rate
	fi
fi

uci commit $CONFIG
uci commit wireless
